# crontab
* * * * * /var/www/html/run-monitor.sh >> /var/www/html/src/Logs/cron.log 2>&1
